# contrib/erlyctl
#!/bin/sh
#
### BEGIN INIT INFO
# Provides:          erlyvideo
# Required-Start:    $network
# Required-Stop:     $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start erlyvideo at boot time and control it in runtime
# Description:       Enable service provided by erlyvideo.
### END INIT INFO
  
  

NODE=ems
HOST=`hostname`
DOTTED=`hostname | grep '\.'`
if [ "" = "$DOTTED" ] ;then
  TYPE="sname"
else
  TYPE="name"
fi


ERL="erl +A 4 +K true"

start()
{
  export ERL_LIBS=deps:lib:plugins:..
  EBIN=`erl -eval 'io:format("~s", [code:lib_dir(erlyvideo,ebin)])' -s init stop -noshell`
  if [ -f "/etc/erlyvideo/production.config" ] ; then
    CONFIG="-config /etc/erlyvideo/production.config"
  else if   [ -f "priv/production.config" ] ; then
    CONFIG="-config priv/production.config"
    else
      CONFIG=""
    fi
  fi
	$ERL -noinput -detached -pa ebin -pa "$EBIN" -pa /var/lib/erlyvideo/plugins $CONFIG \
  -kernel error_logger false \
	-boot start_sasl \
	-s erlyvideo \
	-$TYPE $NODE@$HOST
}

run()
{
  if [ -f "priv/dev.config" ] ; then
    CONFIG="-config priv/dev.config"
  else if [ -f "/etc/erlyvideo/dev.config" ] ; then
    CONFIG="-config /etc/erlyvideo/dev.config"
    else
      CONFIG=""
    fi
  fi
	ERL_LIBS=apps:.. $ERL -pa apps/erlyvideo/ebin -pa /var/lib/erlyvideo/plugins $CONFIG \
	-boot start_sasl \
	-s erlyvideo \
	-$TYPE $NODE@$HOST
	exit
}

shell()
{
    $ERL \
      -$TYPE debug$NODE@$HOST \
      -pa ebin \
      -hidden \
      -remsh $NODE@$HOST
}

ctl()
{
    $ERL \
      -$TYPE ctl$NODE@$HOST \
      -pa ebin \
      -noinput \
      -hidden \
      -s erlyvideo_ctl -extra $NODE@$HOST $@
}

usage()
{
    cat <<-EOF
erlyctl run (in foreground)
erlyctl start (in background)
erlyctl shell (remote shell)
erlyctl stop (shutdown server)
erlyctl reload (validate and reload config file)
erlyctl restart (reload config file and restart whole server)
EOF
    exit
}

stop()
{
  ctl stop
}

reload()
{
  ctl restart
}

restart()
{
  ctl restart
}


[ $# -lt 1 ] && usage

$@

